# lede/rockchip/armv8/base-files/etc/init.d/repair_ro
#!/bin/sh /etc/rc.common

START=22

start() {
    /sbin/repair_ro
}
